# systemd/aerod.service
[Unit]
Description=aerotools daemon

[Service]
ExecStart=/usr/bin/aerod -F -t

[Install]
WantedBy=multi-user.target
